# supervisor.conf
[program:baby_tracker]
directory=/srv/baby-tracker
command=/srv/baby-tracker/main.sh
environment=DB_FILE="./db.sqlite"
stderr_logfile=/var/log/baby_tracker.err.log
stdout_logfile=/var/log/baby_tracker.out.log
stopsignal=INT
stopasgroup=true
killasgroup=true
